# services/hydra/run
#!/bin/bash
HYDRA_RUN_PATH=/hydra
/sbin/setuser abc python $HYDRA_RUN_PATH/nzbhydra.py --nobrowser --logfile /config/nzbhydra.log --database /config/nzbhydra.db --config /config/settings.cfg 
